# services/pcam-overview.timer
[Unit]
Description=Run pcam overview generator every day

[Timer]
OnCalendar=daily
Unit=pcam-overview.service

[Install]
WantedBy=multi-user.target
